# vagrant/oracle/puppet/modules/weblogic/templates/domain_service_startup.erb
#!/bin/bash
# chkconfig: 345 20 80
### BEGIN INIT INFO
# Provides: weblogic
# Required-Start: 
# Default-Start: 3 4 5
# Default-Stop: 0 1 6
# Description: weblogic domain Init Script
### END INIT INFO
PATH=/usr/bin:/sbin:/bin:/usr/sbin
export PATH
RETVAL=0
# Source function library.
. /etc/rc.d/init.d/functions
start() {
	sh <%= @domain_home %>/startWebLogic.sh &
    RETVAL=$?
    return $RETVAL
}
stop() {
	sh <%= @domain_home %>/bin/stopWebLogic.sh &
    RETVAL=$?
    return $RETVAL
}
restart() {
    stop
    start
}
case "$1" in
    start)
        start
    ;;
    stop)
        stop
    ;;
    restart)
        restart
    ;;
    *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 1
esac

exit $RETVAL
